# codex-daemon/codex-daemon.service
[Unit]
Description=Codex Persistent Daemon (OpenAI)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/claudebot/codex-daemon
ExecStart=/usr/bin/python3 %h/claudebot/codex-daemon/daemon.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
