# claude-daemon/claude-daemon.service
[Unit]
Description=Claude Persistent Daemon
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/claudebot/claude-daemon
ExecStart=/usr/bin/python3 %h/claudebot/claude-daemon/daemon.py
Restart=on-failure
RestartSec=10
StandardOutput=append:%h/claudebot/claude-daemon/daemon.log
StandardError=append:%h/claudebot/claude-daemon/daemon.log
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
