# web-chat/claude-web-chat.service
[Unit]
Description=Multi-Agent Web Chat Interface
After=network.target claude-daemon.service

[Service]
Type=simple
WorkingDirectory=%h/claudebot/web-chat
ExecStartPre=/bin/bash -c 'fuser -k 5003/tcp 2>/dev/null || true'
ExecStart=/usr/bin/python3 %h/claudebot/web-chat/app.py
Restart=always
RestartSec=5
Environment=HOME=%h
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
